# speedtest/speedtest.cron
MAILTO=""
*/5 * * * * /root/scripts/speedcheck.sh
