# speedtest/speedtest.cron
MAILTO=""
*/2 * * * * /root/scripts/speedcheck.sh
